# rhel-project/cart.service
[Unit]
Description= Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis.dev.project.online
Environment=CATALOGUE_HOST=catalogue.dev.project.online
ExecStart=/bin/node /app/server.js
SyslogIdentifier=Cart

[Install]
WantedBy=multi.user.target